# cron/mailman.cron
*/2 * * * * * /Users/joe/development/products/grpbldr/env/dev/grpbldr/script/run_mail.sh